# deployment/portfolio-backend.service
[Unit]
Description=Portfolio Backend - Spring Boot API
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user

# Working directory
WorkingDirectory=/home/ec2-user

# Java command to run the JAR
ExecStart=/usr/bin/java -jar /home/ec2-user/portfolio-backend.jar

# Environment variables
Environment="SERVER_PORT=8080"
Environment="SPRING_PROFILES_ACTIVE=production"

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=portfolio-backend

# Resource limits
MemoryMax=300M
MemoryHigh=250M

[Install]
WantedBy=multi-user.target
